# deploy/home-server.service
# Install with:
#   sudo cp deploy/home-server.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now home-server
#
# Adjust User, WorkingDirectory and ExecStart if your username or install
# path differ from the defaults below.

[Unit]
Description=home-server (personal file server)
After=network.target local-fs.target
RequiresMountsFor=/mnt/usb

[Service]
Type=simple
User=sty00a4
Group=sty00a4
WorkingDirectory=/home/sty00a4/home-server
ExecStart=/home/sty00a4/home-server/start.sh
Restart=on-failure
RestartSec=2
Environment=RUST_LOG=info

# --- light sandboxing ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/mnt/usb/ /home/sty00a4/home-server

[Install]
WantedBy=multi-user.target
